# my-ftp-browser/rootfs/etc/service.d/ftp-server/run
#!/usr/bin/with-contenv bashio
# Démarrer le serveur API Python
exec python3 /usr/share/ftpbrowser/api/server.py
