# my-ftp-browser/rootfs/etc/service.d/ftp-server/run
#!/usr/bin/with-contenv bashio
# Démarrer le serveur API Python
bashio::log.info "Starting FTP Browser Share API..."
exec python3 /usr/share/ftpbrowser/api/server.py

