# my-ftp-browser/rootfs/etc/service.d/ftp-server/run
#!/usr/bin/with-contenv bashio
# Démarrer Nginx
exec nginx -g "daemon off;"
