# my-ftp-browser/rootfs/etc/service.d/ftp-server/run
#!/bin/bash
# Service API Python
echo "Démarrage de l'API FTP Browser Share..."
exec python3 /usr/share/ftpbrowser/api/server.py

