# my-ftp-browser/rootfs/etc/service.d/ftp_browser/run
#!/usr/bin/with-contenv bashio
/run.sh

