# my-ftp-browser/rootfs/etc/service.d/ftp_browser/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Démarrage du service FTP Browser"

exec /run.sh

