# my-ftp-browser/rootfs/etc/service.d/nginx/run
#!/usr/bin/with-contenv bashio
# Démarrer Nginx
bashio::log.info "Starting Nginx..."
exec nginx -g "daemon off;"
