# my-ftp-browser/rootfs/etc/service.d/nginx/run
#!/bin/bash
# Service Nginx
echo "Démarrage de Nginx..."
exec nginx -g "daemon off;"

